# crontab
*/2 * * * * /home/ubuntu/back_up.sh /home/ubuntu/test